# services/sabnzbd/run
#!/bin/bash

exec /sbin/setuser abc /usr/bin/sabnzbdplus --config-file /config
